# relysyspackages/Modbus_Package/files/etc/uci-defaults/02_firstboot
#!bin/sh 
. /lib/functions.sh

cp /Web_Page_Gateway_Apps/Common_GW_Files/rc.local /etc

exit 0
